# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable appliancemonitor.service
systemctl start appliancemonitor.service
