# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop appliancemonitor.service
systemctl disable appliancemonitor.service
